# .docker/conf/supervisor_conf/queue-worker.conf
[program:queue-worker]
command=php artisan queue:work
directory=/var/www
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log